# romana-install/roles/romana-common/templates/dnsmasq-calico.conf
# dnsmasq-calico - Customized dnsmasq
#
# dnsmasq-calico description

description	"dnsmasq customizedby calico"

start on net-device-up

respawn
respawn limit 10 5

exec {{ romana_dir }}/bin/dnsmasq-calico --no-daemon --dhcp-range={{ romana[ec2_tag_aws_cloudformation_logical_id].gateway }},{{ romana[ec2_tag_aws_cloudformation_logical_id].gateway }},static,255.255.0.0,infinite --bridge-interface=romana-gw,tap* --interface romana-gw --interface tap* --read-ethers --bind-dynamic
