# romana-install/roles/romana/install-master/templates/upstart/romana-ipam.conf
# romana-ipam - Cloud Native Networking IPAM Service
#
# romana-ipam detailed description

description	"Romana Cloud Native Networking IPAM Service"

start on started romana-root and started romana-topology and started romana-tenant

respawn
respawn limit 10 5

script
    if [ -r {{ romana_etc_dir }}/ipam ]; then
        . {{ romana_etc_dir }}/ipam || true
    fi
    exec {{ romana_bin_dir }}/ipam -rootURL http://{{ stack_nodes.Controller.mgmt_ip }}:9600
end script

