# romana-install/roles/romana/install-master/templates/upstart/romana-policy.conf
# romana-policy - Cloud Native Networking Policy Service
#
# romana-policy detailed description

description	"Romana Cloud Native Networking Policy Service"

start on started romana-root and started romana-topology and started romana-tenant

respawn
respawn limit 10 5

script
    if [ -r {{ romana_etc_dir }}/policy ]; then
        . {{ romana_etc_dir }}/policy || true
    fi
    exec {{ romana_bin_dir }}/policy -rootURL http://{{ romana_master_ip }}:9600
end script
